# ai-lab/service/com.lab01.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lab01.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mike/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>daemon.py</string>
        <string>--daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mike/projects/01/ai-lab</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mike/projects/01/ai-lab/logs/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mike/projects/01/ai-lab/logs/daemon.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mike/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
